# IML.DeviceScannerDaemon/systemd-units/block-device-populator.service
[Unit]
Description=IML Block Device Populator

[Service]
ExecStart=/usr/sbin/udevadm trigger --action=change --subsystem-match=block
Type=oneshot
