# IML.DeviceScannerDaemon/systemd-units/device-scanner.target
[Unit]
Description=Device Scanner Services

Requires=block-device-populator.service
Requires=zed-populator.service
