# IML.ScannerProxyDaemon/systemd-units/scanner-proxy.path
[Unit]
Description=IML Scanner Proxy Path
PartOf=device-scanner.socket
After=device-scanner.socket

[Path]
DirectoryNotEmpty=/var/lib/chroma/settings

[Install]
WantedBy=device-scanner.socket
